# tdiary/build/etc/supervisor/conf.d/tdiary.conf
[program:tdiary]
command = /home/debian/tdiary/tdiary-core/start.sh
numprocs = 1
directory = /home/debian/tdiary/tdiary-core
umask = 022
autostart = true
autorestart = true
startsecs = 3
startretries = 5
user = debian
stdout_logfile = /var/log/supervisor/tdiary.log
stderr_logfile = /var/log/supervisor/tdiary.log
redirect_stderr = true
stopasgroup = true
killasgroup = true

