# ruby/build/etc/profile.d/rbenv.sh
export RBENV_ROOT=/opt/rbenv
export PATH=$RBENV_ROOT/bin:$PATH
eval "$(rbenv init -)"
