# norikra/build/etc/supervisor/conf.d/norikra.conf
[program:norikra]
command = /opt/rbenv/shims/norikra start -s /home/debian/norikra/data/stats.json
numprocs = 1
directory = /home/debian/norikra
umask = 022
autostart = true
autorestart = true
startretries = 1
user = debian
stdout_logfile = /var/log/supervisor/norikra.log
stderr_logfile = /var/log/supervisor/norikra.log
redirect_stderr = true

environment=RBENV_VERSION=jruby-1.7.19,LANG=ja_JP.UTF-8

