# geminabox/build/etc/supervisor/conf.d/geminabox.conf
[program:geminabox]
command = /home/debian/geminabox/start.sh
numprocs = 1
directory = /home/debian/geminabox
umask = 022
autostart = true
autorestart = true
startsecs = 3
startretries = 5
user = debian
stdout_logfile = /var/log/supervisor/geminabox.log
stderr_logfile = /var/log/supervisor/geminabox.log
redirect_stderr = true
stopasgroup = true
killasgroup = true

