# nodejs/build/etc/profile.d/ndenv.sh
export NDENV_ROOT=/opt/ndenv
export PATH=$NDENV_ROOT/bin:$PATH
eval "$(ndenv init -)"
