# nodejs/build/etc/profile.d/nvm.sh
if [ -r /opt/nvm/nvm.sh ]; then
  . /opt/nvm/nvm.sh
fi
