# fluentd-ui/build/etc/supervisor/conf.d/fluentd-ui.conf
[program:fluentd-ui]
command = /opt/fluentd-ui/start_server.sh
numprocs = 1
directory = /opt/fluentd-ui
umask = 022
autostart = true
autorestart = true
startretries = 1
user = td-agent
stdout_logfile = /var/log/supervisor/fluentd-ui.log
stderr_logfile = /var/log/supervisor/fluentd-ui.log
redirect_stderr = true
stopasgroup = true
killasgroup = true

environment=HOME=/home/td-agent,PATH=/usr/local/bin:/usr/bin:/bin,LANG=ja_JP.UTF-8

