# etcd/build/etc/supervisor/conf.d/etcd.conf
[program:etcd]
command = bash -l -c "exec etcd"
numprocs = 1
directory = /var/lib/etcd
umask = 022
autostart = true
autorestart = true
startretries = 1
user = etcd
stdout_logfile = /var/log/supervisor/etcd.log
stderr_logfile = /var/log/supervisor/etcd.log
redirect_stderr = true
environment=HOME="/var/lib/etcd"
