# gemstash/build/etc/supervisor/conf.d/gemstash.conf
[program:gemstash]
command = gemstash start --no-daemonize
numprocs = 1
directory = /home/debian/
umask = 022
autostart = true
autorestart = true
startsecs = 3
startretries = 5
user = debian
stdout_logfile = /var/log/supervisor/gemstash.log
stderr_logfile = /var/log/supervisor/gemstash.log
redirect_stderr = true
stopasgroup = true
killasgroup = true
environment = USER="debian",HOME="/home/debian",LANG="en_US.UTF-8",RACK_ENV="production"

